# scripts/com.prevoir.poll-jira.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd schedule for poll-jira.sh
  Runs at the top of every hour (StartCalendarInterval); reliable after sleep/wake.

  SETUP:
  1. Replace /Users/YOUR_USERNAME below with your actual home directory path
  2. Copy this file to ~/Library/LaunchAgents/
  3. Run: launchctl load ~/Library/LaunchAgents/com.prevoir.poll-jira.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prevoir.poll-jira</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- UPDATE: replace /Users/YOUR_USERNAME with your actual home path -->
        <string>/Users/YOUR_USERNAME/Documents/Prevoir/Scripts/poll-jira.sh</string>
    </array>

    <!-- Run at the top of every hour; more reliable than StartInterval after sleep -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- UPDATE: replace /Users/YOUR_USERNAME with your actual home path -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Documents/Prevoir/Scripts/poll-jira.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Documents/Prevoir/Scripts/poll-jira-error.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
